# examples/lxc/ROOTFS/etc/anopa/services/gettycon/run
#!/command/execlineb -P
getty -L -H CONSOLE 38400 console vt100
