# examples/lxc/ROOTFS/etc/anopa/services/dropbear/run
#!/command/execlineb -P
/command/dropbear -F -E
